# all/etc/systemd/system/pprotein.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/pprotein

[Install]
WantedBy=multi-user.target